# .zsh/profile.d/xmodmap.zsh
which xmodmap &> /dev/null && \
xmodmap - << EOF
remove Lock = Super_L
keysym Super_L = F12
EOF
