# scripts/systemd/user-hook-shutdown@.service
[Unit]
Description=User's shutdown hook
ConditionPathExists=/home/%I/.hooks/shutdown.sh
After=multi-user.target

[Service]
User=%I
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/true
ExecStop=/bin/sh -e /home/%I/.hooks/shutdown.sh

[Install]
WantedBy=multi-user.target
